# python_backend/crypto_ml_trading/systemd/crypto-ml-trading.service
[Unit]
Description=Crypto ML Trading System
Documentation=https://github.com/yourcompany/crypto-ml-trading
After=network.target
Wants=network.target

[Service]
Type=simple
User=trading
Group=trading
WorkingDirectory=/opt/crypto-ml-trading
ExecStart=/usr/bin/python3 -u main.py
ExecReload=/bin/kill -HUP $MAINPID

# Environment variables
Environment=ENVIRONMENT=production
Environment=LOG_LEVEL=INFO
Environment=PYTHONPATH=/opt/crypto-ml-trading
Environment=PYTHONUNBUFFERED=1

# Restart policy
Restart=on-failure
RestartSec=5
StartLimitInterval=60s
StartLimitBurst=3

# Resource limits
MemoryLimit=4G
CPUQuota=200%

# Security settings
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/opt/crypto-ml-trading/data
ReadWritePaths=/opt/crypto-ml-trading/logs
ReadWritePaths=/opt/crypto-ml-trading/saved_models

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=crypto-ml-trading

[Install]
WantedBy=multi-user.target